# teacher/plc/otlab-plc.service
[Unit]
Description=OTLab ladder-logic PLC engine (scan loop driving Qwiic I/O)
Documentation=https://github.com/aaroncorvo/otlab
# Order after whichever :8090 backend is enabled (physical otlab-qwiic OR
# virtual otlab-vio) but do NOT hard-require either — the engine tolerates
# the backend being down (it retries the REST API), and a hard Requires
# would resurrect otlab-qwiic when switching to the virtual backend.
After=network-online.target otlab-qwiic.service otlab-vio.service
Wants=network-online.target

[Service]
Type=simple
Environment=OTLAB_PLC_PORT=8091
Environment=OTLAB_QWIIC_URL=http://127.0.0.1:8090
Environment=OTLAB_PLC_STORE=/var/lib/otlab/plc/program.json
Environment=DASH_USER=otlab
Environment=DASH_PASS=P@ssw0rd!
ExecStart=/usr/bin/python3 /opt/otlab/otlab_plc.py
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
